# systemd/k_e_r-daily.timer
[Unit]
Description=K_E_R daily note — 매일 16:00 KST (장 마감 +30분)
Documentation=https://github.com/soccz/K_E_R

[Timer]
# 평일(월~금) 16:00 KST (KRX 장 마감 15:30 + 30분 버퍼).
# 토·일은 KRX 휴장 → 데이터 갱신 없음 → 검사 자체 skip.
# 한국 공휴일도 KRX 휴장이지만 systemd는 모름 — trigger 모듈이 OHLCV 변동 0이면 자동 skip.
OnCalendar=Mon..Fri *-*-* 16:00:00
# 슬롯 시간에 컴퓨터가 꺼져있었으면 부팅 즉시 catch-up (휴일은 거래일 데이터로 한 번 검사)
Persistent=true
RandomizedDelaySec=120
AccuracySec=1min

[Install]
WantedBy=timers.target
