# systemd/k_e_r-weekend.service
[Unit]
Description=K_E_R weekend industry note
Documentation=https://github.com/soccz/K_E_R
After=network-online.target
Wants=network-online.target
OnFailure=k_e_r-notify-failure@%n.service

[Service]
Type=oneshot
WorkingDirectory=/home/soccz/22tb/report
EnvironmentFile=/home/soccz/22tb/report/.env
ExecStart=/home/soccz/22tb/report/bin/run_weekend.sh
StandardOutput=append:/home/soccz/22tb/report/logs/weekend.log
StandardError=append:/home/soccz/22tb/report/logs/weekend.log
TimeoutStartSec=3600
Nice=10
# exit 3 = LLM 가드 차단 (rate limit 등 정상 skip). 실패 아니므로 success 처리.
SuccessExitStatus=3

[Install]
WantedBy=default.target
